# scripts/com.ima.sellmodel.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ima.sellmodel.daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/amardani/Sell-Model/scripts/daily_refresh.sh</string>
  </array>
  <!-- Weekday mornings at 06:45 local time. If the Mac is asleep at 06:45,
       launchd runs the job at the next wake; if it is shut down, the day is
       skipped and the site keeps the previous data. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>45</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>6</integer><key>Minute</key><integer>45</integer></dict>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/amardani/.npm-global/bin:/Users/amardani/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/amardani</string>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/amardani/Sell-Model/output/daily/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/amardani/Sell-Model/output/daily/launchd.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
